# download-sweeper.timer
[Unit]
Description=Run download-sweeper daily and on boot

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
